# mediawiki/etc/crontab
# min   hour    day     month   weekday command
*/3     *       *       *       *       /usr/bin/php /var/www/mediawiki/maintenance/runJobs.php 2>&1
