# mediawiki/etc/crontab
# min   hour    day     month   weekday command
*/2     *       *       *       *       /usr/bin/php /var/www/mediawiki/maintenance/runJobs.php --quiet 2>&1
